# root/etc/supervisor/conf.d/bind.conf
[program:named]
command=/usr/sbin/named -u named -f
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
